# crontab
* * * * * flock -n /tmp/shutdown_script.lock /usr/src/app/run_script.sh >> /usr/src/app/logs/cron-logs.log 2>&1
